# anode-system/scraper/anode-scheduler.service
[Unit]
Description=Anode Scraper Scheduler Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/cost-calculator/anode-system/scraper
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /path/to/cost-calculator/anode-system/scraper/scheduler.py
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/anode-scheduler.log
StandardError=append:/var/log/anode-scheduler-error.log

[Install]
WantedBy=multi-user.target